# packaging/debian/ferrite-tray.service
[Unit]
Description=Ferrite tray (Android-as-external-display)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ferrite-tray
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
